# xyz.openbmc_project.nf.power.manager.service.in
[Unit]
Description=NF card Power Control

[Service]
Restart=always
RestartSec=3
StartLimitBurst=10
ExecStart=@bindir@/nf_pwr_ctrl
SyslogIdentifier=nf_pwr_ctrl
Type=dbus
BusName=xyz.openbmc_project.nf.power.manager

[Install]
WantedBy=multi-user.target
